# backend/server_setup/gunicorn.socket
[Unit]
Description=Threes gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
